# systemd/taptunes.service
[Unit]
Description=TapTunes Music System
Documentation=https://github.com/yourusername/TapTunes
After=network.target sound.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/greg/taptunes/python-services
ExecStart=/home/greg/taptunes/venv/bin/python /home/greg/taptunes/python-services/taptunes_main.py
Restart=always
RestartSec=10

# Environment variables
Environment="NODE_ENV=production"
Environment="BACKEND_URL=http://localhost:3001"
Environment="TAPTUNES_INSTALL_DIR=/home/greg/taptunes"
Environment="PATH=/home/greg/taptunes/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Resource limits
LimitNOFILE=4096
MemoryMax=512M
CPUQuota=95%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taptunes

[Install]
WantedBy=multi-user.target
