# scripts/taptunes-rfid.service
[Unit]
Description=TapTunes RFID Scanner
After=network.target taptunes-backend.service

[Service]
Type=simple
User=greg
WorkingDirectory=/home/greg/TapTunes/rfid
ExecStart=/usr/bin/python3 rfid_scanner.py
Restart=always
RestartSec=10

# System limits for Pi Zero W
MemoryMax=64M
CPUQuota=20%

[Install]
WantedBy=multi-user.target