# scripts/taptunes-backend.service
[Unit]
Description=TapTunes Backend
After=network.target

[Service]
Type=simple
User=greg
WorkingDirectory=/home/greg/TapTunes/backend
Environment="NODE_ENV=production"
Environment="PORT=3001"
Environment="HOST=0.0.0.0"
ExecStart=/usr/local/bin/node dist/index.js
Restart=always
RestartSec=10

# System limits for Pi Zero W
LimitNOFILE=4096
MemoryMax=256M
CPUQuota=80%

[Install]
WantedBy=multi-user.target